# Dockerfile.e2e.supervisor.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:postgresql]
command=/usr/bin/postgres -D /var/lib/postgresql/data
user=postgres
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql_err.log

[program:api]
command=pnpm --filter @codecomply/api start
directory=/app
autostart=false
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api_err.log
environment=PORT=4000,DATABASE_URL="postgresql://inspector:test_password@localhost:5432/inspections_test",NODE_ENV=test,JWT_SECRET=test_jwt_secret

[program:inspector]
command=pnpm --filter @codecomply/inspector preview --port 5175 --host 0.0.0.0
directory=/app
autostart=false
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/inspector.log
stderr_logfile=/var/log/supervisor/inspector_err.log

[program:admin]
command=pnpm --filter @codecomply/admin preview --port 5174 --host 0.0.0.0
directory=/app
autostart=false
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/admin.log
stderr_logfile=/var/log/supervisor/admin_err.log
